# scripts/com.agent-platform.gateway.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agent-platform.gateway</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON3_PATH__</string>
        <string>-m</string>
        <string>uvicorn</string>
        <string>gateway.main:app</string>
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>8000</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PROJECT_DIR__</string>

    <!-- 开机自动启动 -->
    <key>RunAtLoad</key>
    <true/>

    <!-- 崩溃后自动重启，最快间隔 10 秒 -->
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <!-- 日志输出 -->
    <key>StandardOutPath</key>
    <string>__PROJECT_DIR__/logs/launchd_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__PROJECT_DIR__/logs/launchd_stderr.log</string>

    <!-- 环境变量 -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
